# services/hamachi/supervisord_hamachi.conf
;=======================================
; Hamachi service
;=======================================
 
[program:hamachi]
 
; General
;command       = /etc/supervisor/conf.d/run_Hamachi.sh
command       = /opt/logmein-hamachi/bin/hamachid -c /data/hamachi_config debug
user          = root
numprocs      = 1
autostart     = true
autorestart   = true
startsecs     = 10
stopwaitsecs  = 30
process_name  = hamachi
 
; Standard out / error
stdout_logfile          = /var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups  = 10
stderr_logfile          = /var/log/supervisor/%(program_name)s.log
stderr_logfile_maxbytes = 5MB
stderr_logfile_backups  = 10
